# installation/mercure_receiver.service
[Unit]
Description=mercure Receiver
Documentation=https://mercure-imaging.org/docs
After=network.target mercure_bookkeeper.service
Requires=network.target
Wants=mercure_bookkeeper.service

[Service]
Type=simple
WorkingDirectory=/opt/mercure/app
ExecStart=/opt/mercure/app/receiver.sh
Restart=on-failure
RestartSec=3
User=mercure
EnvironmentFile=/opt/mercure/config/receiver.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mercure/data /opt/mercure/config
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
