# installation/mercure_processor.service
[Unit]
Description=mercure Processor
Documentation=https://mercure-imaging.org/docs
After=network.target mercure_bookkeeper.service
Requires=network.target
Wants=mercure_bookkeeper.service

[Service]
Type=simple
WorkingDirectory=/opt/mercure/app
ExecStart=/opt/mercure/env/bin/python /opt/mercure/app/processor.py
Restart=on-failure
RestartSec=3
User=mercure
SupplementaryGroups=docker

# Security hardening (limited - needs Docker access)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mercure/data /opt/mercure/config /var/run/docker.sock
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictRealtime=true
RestrictSUIDSGID=true
# Note: Cannot use RestrictNamespaces or ProtectControlGroups - needs Docker access

[Install]
WantedBy=multi-user.target
